# i_cant_computer/vanilla_rhel_deploy/roles/automate_backups/templates/backup-agent.timer.j2
[Unit]
Description=Run backups daily

[Timer]
OnCalendar={{ automate_backups_schedule }}

[Install]
WantedBy=timers.target
